# deuxiemeDockerfile/supervisor.conf
[supervisord]
nodaemon=true

[program:apache2]
priority=10
directory=/tmp
command=/usr/sbin/apache2ctl -D FOREGROUND
user=root
autostart=true
autorestart=true

[program:mariadb]
command=/etc/init.d/mariadb start
priority=999
